# roles/ctlabs_k3s/templates/k3s.service.j2
[Unit]
Description=k3s
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple

User=root

ExecStart={{ ctlabs_k3s.defaults.pkgs.k3s.path }} server

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
