# roles/ctlabs_openclaw/templates/openclaw.service.j2
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/openclaw gateway \
  --bind lan \
  --port 18789 \
  --allow-unconfigured \
  --tls-enabled \
  --tls-cert-path /etc/ca-ctlabs/h1.ctlabs.internal.crt \
  --tls-key-path /etc/ca-ctlabs/h1.ctlabs.internal.prv \
  --tls-min-version TLSv1.2 \
  --tls-ciphers "HIGH:!aNULL:!MD5" \
  --control-ui-allowed-origins "http://192.168.100.166:18789,https://192.168.100.166:18789"

Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group

[Install]
WantedBy=default.target
