# roles/ctlabs_openclaw/templates/openclaw.service.j2
[Unit]
Description=OpenClaw Gateway
After=network-online.target
Wants=network-online.target

[Service]
User={{ ctlabs_openclaw.defaults.service.user }}
Group={{ ctlabs_openclaw.defaults.service.group }}
{% if ctg_os.startswith('kali') %}
ExecStart=/usr/local/bin/openclaw gateway --bind lan --port 18789 --allow-unconfigured 
{% else %}
ExecStart=/usr/bin/openclaw gateway --bind lan --port 18789 --allow-unconfigured 
{% endif %}


Restart=always
RestartSec=5
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
KillMode=control-group

[Install]
WantedBy=default.target
